# Test_FW/MIPS/_V1.00(AARO.1)C0.bin-0.extracted/jffs2-root/fs_1/usr/lib/opkg/info/app-agent.list
/etc/init.d/app-agent
/usr/include/app_notify.h
/usr/lib/libmini.so
/usr/bin/EmappC
/usr/bin/EmappS
